# autoinstall/codis/conf/dashboard-supervisord.conf
[unix_http_server]
file=INSTALL_DIR/run/codis/dashboard/supervisor.sock

[supervisord]
logfile=INSTALL_DIR/logs/codis/dashboard/supervisord.log
logfile_maxbytes=1024MB
logfile_backups=10
loglevel=info
pidfile=INSTALL_DIR/run/codis/dashboard/supervisord.pid
nodaemon=false
minfds=10240                ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://INSTALL_DIR/run/codis/dashboard/supervisor.sock

[program:codis-dashboard]
command=INSTALL_DIR/codis/bin/codis-config -c INSTALL_DIR/codis/conf/config.ini -L INSTALL_DIR/logs/codis/dashboard/dashboard.log dashboard --addr=:18087 --http-log=INSTALL_DIR/logs/codis/dashboard/requests.log
redirect_stderr=true
stdout_logfile=INSTALL_DIR/logs/codis/dashboard/%(program_name)s.console.log
autorestart=true

