# autoinstall/codis/conf/proxy-supervisord.conf
[unix_http_server]
file=INSTALL_DIR/run/codis/proxy/supervisor.sock

[supervisord]
logfile=INSTALL_DIR/logs/codis/proxy/supervisord.log
logfile_maxbytes=1024MB
logfile_backups=10
loglevel=info
pidfile=INSTALL_DIR/run/codis/proxy/supervisord.pid
nodaemon=false
minfds=10240                ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://INSTALL_DIR/run/codis/proxy/supervisor.sock

[program:codis-proxy]
command=INSTALL_DIR/codis/bin/codis-proxy --log-level warn -c INSTALL_DIR/codis/conf/config.ini -L INSTALL_DIR/logs/codis/proxy/proxy.log  --cpu=8 --addr=PROXY_IPADDR:16379 --http-addr=PROXY_IPADDR:11000
redirect_stderr=true
stdout_logfile=INSTALL_DIR/logs/codis/proxy/%(program_name)s.console.log
autorestart=true

